# config/supervisor.conf
[supervisord]
nodaemon=true

[program:hiveserver2]
command=/usr/local/hive/bin/hiveserver2
